# LmDbServer/scripts/deleteObsolete.cron
#!/bin/bash

# @summary: This script will delete old occurrencesets cascading to dependent objects

COUNTER=0
while [ $COUNTER -lt 10 ]; do
   date
   echo Counter is $COUNTER 
   psql -U admin -d mal -c "DELETE FROM occurrenceset WHERE occurrencesetid IN 
                             (SELECT occurrencesetid FROM occurrenceset 
                                WHERE userid = 'lm2' AND scientificnameid IS NULL 
                                LIMIT 10);"
   let COUNTER=COUNTER+1
   sleep 10
done
