# LmDbServer/scripts/deleteOldAnonymous.cron.in
#!/bin/bash

# @author: CJ Grady
# @lastModified: 2013/04/30 16:59
# @summary: This script will delete old anonymous data

# Set environment in case it is not already set
. /etc/profile.d/lifemapper.sh

# Have to export PYTHONPATH or the modules will not be found
export PYTHONPATH=@LMHOME@

@PYBIN@ @LMHOME@/LmDbServer/scripts/deleteOldAnonymous.py
