# LmWebServer/scripts/cron/buildLmLuceneIndex.cron.in
#!/bin/bash

# @author: CJ Grady
# @lastModified: 2015/10/21 12:37
# @summary: This script will rebuild the Lifemapper species lucene index

# Set environment in case it is not already set
. /etc/profile.d/lifemapper.sh

# Have to export PYTHONPATH or the modules will not be found
#   CJG (10/21/2015): I don't think that this is true now, test this
export PYTHONPATH=@LMHOME@

# Build the index
# Note that the output of this could be piped to a file to check for success
#   or you can check the last modified time of the files in: 
#   @LMWEB@/luceneIndex/

@PYBIN@ @LMHOME@/LmWebServer/lucene/lmLucene.py build
