# LmDebug/tests/end2end/scripts/end2endTests.cron.in
#!/bin/bash

# @author: CJ Grady
# @lastModified: 2014/08/27 12:04
# @summary: This script will run Lifemapper end to end tests
# @license: gpl2
# @copyright: Copyright (C) 2013, University of Kansas Center for Research
#  
#            Lifemapper Project, lifemapper [at] ku [dot] edu, 
#            Biodiversity Institute,
#            1345 Jayhawk Boulevard, Lawrence, Kansas, 66045, USA
#     
#            This program is free software; you can redistribute it and/or modify 
#            it under the terms of the GNU General Public License as published by 
#            the Free Software Foundation; either version 2 of the License, or (at 
#            your option) any later version.
#    
#            This program is distributed in the hope that it will be useful, but 
#            WITHOUT ANY WARRANTY; without even the implied warranty of 
#            MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the GNU 
#            General Public License for more details.
#    
#            You should have received a copy of the GNU General Public License 
#            along with this program; if not, write to the Free Software 
#            Foundation, Inc., 51 Franklin Street, Fifth Floor, Boston, MA 
#            02110-1301, USA.

# Set environment in case it is not already set
. /etc/profile.d/lifemapper.sh

# Nadya: this should read the configured APP_DIR
# Have to export PYTHONPATH or the modules will not be found
export PYTHONPATH=@LMHOME@

@PYBIN@ @LMHOME@/test/end2end/rad/experiment.py

@PYBIN@ @LMHOME@/test/end2end/sdm/experiment.py

